# systemd/mapper.service.in
[Unit]
Description=WebSocket Webcam Mapper

[Service]
Type=simple
User=@USER@
Group=@GROUP@
ExecStart=@PYTHON@ @MAPPER@ @HOST@ @PORT@

[Install]
WantedBy=multi-user.target
