# systemd/recorder.service.in
[Unit]
Description=WebSocket Webcam Recorder
Requires=@PREFIX@mapper.service
After=@PREFIX@mapper.service

[Service]
Type=simple
User=@USER@
Group=@GROUP@
ExecStart=@PYTHON@ @RECORDER@ @HOST@ @PORT@

[Install]
WantedBy=multi-user.target
