# systemd/server.service.in
[Unit]
Description=WebSocket Webcam Server
Requires=@PREFIX@mapper.service
Requires=@PREFIX@recorder.service
After=@PREFIX@mapper.service
After=@PREFIX@recorder.service

[Service]
Type=simple
User=@USER@
Group=@GROUP@
ExecStart=@PYTHON@ @SERVER@ @HOST@ @PORT@

[Install]
WantedBy=multi-user.target
