# systemd/server.service.in
[Unit]
Description=WebSocket Webcam Server
Requires=@MAPPER_SERVICE@
Requires=@RECORDER_SERVICE@
After=@MAPPER_SERVICE@
After=@RECORDER_SERVICE@

[Service]
Type=simple
User=@USER@
Group=@GROUP@
ExecStart=@PYTHON@ @SERVER@ @PORT@

[Install]
WantedBy=multi-user.target
